# resource/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
nodaemon = true
logfile  = /tmp/supervisord.log

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:mysql]
user        = root
command     = /usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
autostart   = true
autorestart = true

[program:elasticmq]
user        = root
command     = /usr/bin/java -Dconfig.file=/etc/elasticmq.conf -jar /opt/elasticmq.jar
autostart   = true
autorestart = true

[program:dynamodb]
user        = root
command     = /usr/bin/java -Djava.library.path=/opt/dynamodb_local/DynamoDBLocal_lib -jar /opt/dynamodb_local/DynamoDBLocal.jar --port 8000
autostart   = true
autorestart = true
